# tedge-utils/etc/cron.d/boot_crontab
@reboot root sleep 10 && /usr/lib/tedge-utils/bootup.sh
